# etc/udev/rules.d/70-snap.ffmpeg.rules
# This file is automatically generated.
# opengl
KERNEL=="nvhost-*", TAG+="snap_ffmpeg_ffmpeg"
# opengl
KERNEL=="nvmap", TAG+="snap_ffmpeg_ffmpeg"
# opengl
KERNEL=="renderD[0-9]*", TAG+="snap_ffmpeg_ffmpeg"
# opengl
KERNEL=="tegra_dc_[0-9]*", TAG+="snap_ffmpeg_ffmpeg"
# opengl
KERNEL=="tegra_dc_ctrl", TAG+="snap_ffmpeg_ffmpeg"
# opengl
KERNEL=="vchiq", TAG+="snap_ffmpeg_ffmpeg"
# opengl
SUBSYSTEM=="drm", KERNEL=="card[0-9]*", TAG+="snap_ffmpeg_ffmpeg"
# optical-drive
KERNEL=="scd[0-9]*", TAG+="snap_ffmpeg_ffmpeg"
# optical-drive
KERNEL=="sr[0-9]*", TAG+="snap_ffmpeg_ffmpeg"
# optical-drive
SUBSYSTEM=="scsi_generic", SUBSYSTEMS=="scsi", ATTRS{type}=="4|5", TAG+="snap_ffmpeg_ffmpeg"
TAG=="snap_ffmpeg_ffmpeg", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_ffmpeg_ffmpeg $devpath $major:$minor"
# opengl
KERNEL=="nvhost-*", TAG+="snap_ffmpeg_ffplay"
# opengl
KERNEL=="nvmap", TAG+="snap_ffmpeg_ffplay"
# opengl
KERNEL=="renderD[0-9]*", TAG+="snap_ffmpeg_ffplay"
# opengl
KERNEL=="tegra_dc_[0-9]*", TAG+="snap_ffmpeg_ffplay"
# opengl
KERNEL=="tegra_dc_ctrl", TAG+="snap_ffmpeg_ffplay"
# opengl
KERNEL=="vchiq", TAG+="snap_ffmpeg_ffplay"
# opengl
SUBSYSTEM=="drm", KERNEL=="card[0-9]*", TAG+="snap_ffmpeg_ffplay"
# optical-drive
KERNEL=="scd[0-9]*", TAG+="snap_ffmpeg_ffplay"
# optical-drive
KERNEL=="sr[0-9]*", TAG+="snap_ffmpeg_ffplay"
# optical-drive
SUBSYSTEM=="scsi_generic", SUBSYSTEMS=="scsi", ATTRS{type}=="4|5", TAG+="snap_ffmpeg_ffplay"
TAG=="snap_ffmpeg_ffplay", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_ffmpeg_ffplay $devpath $major:$minor"
# opengl
KERNEL=="nvhost-*", TAG+="snap_ffmpeg_ffprobe"
# opengl
KERNEL=="nvmap", TAG+="snap_ffmpeg_ffprobe"
# opengl
KERNEL=="renderD[0-9]*", TAG+="snap_ffmpeg_ffprobe"
# opengl
KERNEL=="tegra_dc_[0-9]*", TAG+="snap_ffmpeg_ffprobe"
# opengl
KERNEL=="tegra_dc_ctrl", TAG+="snap_ffmpeg_ffprobe"
# opengl
KERNEL=="vchiq", TAG+="snap_ffmpeg_ffprobe"
# opengl
SUBSYSTEM=="drm", KERNEL=="card[0-9]*", TAG+="snap_ffmpeg_ffprobe"
# optical-drive
KERNEL=="scd[0-9]*", TAG+="snap_ffmpeg_ffprobe"
# optical-drive
KERNEL=="sr[0-9]*", TAG+="snap_ffmpeg_ffprobe"
# optical-drive
SUBSYSTEM=="scsi_generic", SUBSYSTEMS=="scsi", ATTRS{type}=="4|5", TAG+="snap_ffmpeg_ffprobe"
TAG=="snap_ffmpeg_ffprobe", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_ffmpeg_ffprobe $devpath $major:$minor"
