# etc/udev/rules.d/70-snap.pulseaudio.rules
# This file is automatically generated.
# pulseaudio
KERNEL=="controlC[0-9]*", TAG+="snap_pulseaudio_pulseaudio"
# pulseaudio
KERNEL=="pcmC[0-9]*D[0-9]*[cp]", TAG+="snap_pulseaudio_pulseaudio"
# pulseaudio
KERNEL=="timer", TAG+="snap_pulseaudio_pulseaudio"
TAG=="snap_pulseaudio_pulseaudio", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_pulseaudio_pulseaudio $devpath $major:$minor"
