# etc/udev/rules.d/70-snap.mpv.rules
# This file is automatically generated.
# opengl
KERNEL=="nvhost-*", TAG+="snap_mpv_mpv"
# opengl
KERNEL=="nvmap", TAG+="snap_mpv_mpv"
# opengl
KERNEL=="renderD[0-9]*", TAG+="snap_mpv_mpv"
# opengl
KERNEL=="tegra_dc_[0-9]*", TAG+="snap_mpv_mpv"
# opengl
KERNEL=="tegra_dc_ctrl", TAG+="snap_mpv_mpv"
# opengl
KERNEL=="vchiq", TAG+="snap_mpv_mpv"
# opengl
SUBSYSTEM=="drm", KERNEL=="card[0-9]*", TAG+="snap_mpv_mpv"
# optical-drive
KERNEL=="scd[0-9]*", TAG+="snap_mpv_mpv"
# optical-drive
KERNEL=="sr[0-9]*", TAG+="snap_mpv_mpv"
# optical-drive
SUBSYSTEM=="scsi_generic", SUBSYSTEMS=="scsi", ATTRS{type}=="4|5", TAG+="snap_mpv_mpv"
TAG=="snap_mpv_mpv", RUN+="/usr/lib/snapd/snap-device-helper $env{ACTION} snap_mpv_mpv $devpath $major:$minor"
